# labsite/setup/supervisord.ini

[program:gunicorn]
command = %(HOME_DIR)svenv/bin/gunicorn_django -w 4 -b 127.0.0.1:8443 --access-logfile=%(HOME_DIR)slog/gunicorn_access.log --error-logfile=%(HOME_DIR)slog/gunicorn_error.log
directory = %(PROJ_DIR)s
autostart = true
autorestart = true
user = %(USERNAME)s
group = %(USERNAME)s
environment=PYTHONPATH='%(HOME_DIR)s.env/bin'

[program:celeryd]
command = %(HOME_DIR)svenv/bin/python manage.py celeryd
stdout_logfile = %(HOME_DIR)slog/celeryd.log
stderr_logfile = %(HOME_DIR)slog/celeryd.log
directory = %(PROJ_DIR)s
autostart = true
autorestart = true
startsecs = 10
stopwaitsecs = 600
user = %(USERNAME)s
group = %(USERNAME)s
environment=PYTHONPATH='%(HOME_DIR)svenv/bin'

[program:celerybeat]
command = %(HOME_DIR)svenv/bin/python manage.py celerybeat
stdout_logfile = %(HOME_DIR)slog/celerybeat.log
stderr_logfile = %(HOME_DIR)slog/celerybeat.log
directory = %(PROJ_DIR)s
autostart = true
autorestart = true
startsecs = 10
stopwaitsecs = 600
user = %(USERNAME)s
group = %(USERNAME)s
environment=PYTHONPATH='%(HOME_DIR)svenv/bin'
