# labsite/setup/supervisord.ini

[program:gunicorn]
command = %(HOME)s/venv/bin/gunicorn_django -w 4 -b 127.0.0.1:8000 --access-logfile=%(HOME)s/log/gunicorn_access.log --error-logfile=%(HOME)s/log/gunicorn_error.log %(HOME)s/labsite/settings.py
autostart = true
autorestart = true
user = labuser
group = labuser
environment=PYTHONPATH='%(HOME)s/.env/bin',DJANGO_SETTINGS_MODULE=settings

[program:celeryd]
stdout_logfile = %(HOME)s/log/celeryd.log
stderr_logfile = %(HOME)s/log/celeryd.log
directory = %(HOME)s/labsite/
autostart = true
autorestart = true
startsecs = 10
stopwaitsecs = 600
user = labuser
group = labuser
environment=PYTHONPATH='%(HOME)s/venv/bin',SETTINGS_MODULE=settings
command = %(HOME)s/venv/bin/python %(HOME)s/manage.py celeryd

[program:celerybeat]
stdout_logfile = %(HOME)s/log/celerybeat.log
stderr_logfile = %(HOME)s/log/celerybeat.log
directory = %(HOME)s/labsite/
autostart = true
autorestart = true
startsecs = 10
stopwaitsecs = 600
user = labuser
group = labuser
environment=PYTHONPATH='%(HOME)s/venv/bin',SETTINGS_MODULE=settings
command = %(HOME)s/venv/bin/python %(HOME)s/manage.py celerybeat
