# deploy/config_examples/supervisord.ini

[program:gunicorn]
command = /opt/lab/labsite/.env/bin/gunicorn_django -w 4 -b 127.0.0.1:8000 --access-logfile=opt/lab/labsite/log/gunicorn_access.log --error-logfile=opt/lab/labsite/log/gunicorn_error.log opt/lab/labsite/labsite/settings.py
autostart = true
autorestart = true
user = labuser
group = labuser
environment=PYTHONPATH='/opt/lab/labsite/.env/bin',DJANGO_SETTINGS_MODULE=settings

[program:celeryd]
stdout_logfile = /opt/lab/labsite/log/celeryd.log
stderr_logfile = /opt/lab/labsite/log/celeryd.log
directory = /opt/lab/labsite/labsite/
autostart = true
autorestart = true
startsecs = 10
stopwaitsecs = 600
user = labuser
group = labuser
environment=PYTHONPATH='/opt/lab/labsite/.env/bin',SETTINGS_MODULE=settings
command = /opt/lab/labsite/.env/bin/python /opt/lab/labsite/manage.py celeryd

[program:celerybeat]
stdout_logfile = /opt/lab/labsite/log/celerybeat.log
stderr_logfile = /opt/lab/labsite/log/celerybeat.log
directory = /opt/lab/labsite/labsite/
autostart = true
autorestart = true
startsecs = 10
stopwaitsecs = 600
user = labuser
group = labuser
environment=PYTHONPATH='/opt/lab/labsite/.env/bin',SETTINGS_MODULE=settings
command = /opt/lab/labsite/.env/bin/python /opt/lab/labsite/manage.py celerybeat
